# eclus.socket
[Unit]
Description=EPMD Server Activation Socket

[Socket]
ListenStream=4369

[Install]
WantedBy=sockets.target
